# configuration/courselore.service
[Unit]
Description=Courselore
After=network.target

[Service]
ExecStart=/root/courselore/courselore /root/courselore/configuration.js
User=root
Restart=always

[Install]
WantedBy=multi-user.target
